# contrib/systemd/village-analyze.service
[Unit]
Description=Village Daily Analysis
Documentation=https://github.com/anomalyco/village
After=network-online.target

[Service]
Type=oneshot
# Replace %h with home directory at install time
WorkingDirectory=%h/source/village
ExecStart=%h/.local/bin/village doctor --json --output %h/.village/analysis-report.json
StandardOutput=journal
StandardError=journal
SyslogIdentifier=village-doctor
TimeoutStartSec=300

[Install]
WantedBy=default.target
